# scripts/uvicorn.health-status.path
[Path]
PathChanged=/opt/gitrepo/scripts/health-status.py

[Install]
WantedBy=multi-user.target